# wpa_supplicant/systemd/wpa_supplicant.service.in
[Unit]
Description=WPA supplicant

[Service]
Type=dbus
BusName=fi.epitest.hostap.WPASupplicant
ExecStart=@BINDIR@/wpa_supplicant -u

[Install]
WantedBy=multi-user.target
Alias=dbus-fi.epitest.hostap.WPASupplicant.service
